# ops/auto-parse-65.service
[Unit]
Description=Auto Parse Next.js production service on 65
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/auto-parse
Environment=NODE_ENV=production
Environment=PORT=1007
Environment=HOSTNAME=127.0.0.1
Environment=ADSPOWER_DISABLE_AUTO_START=true
Environment=ADS_POOL_ENABLE_CDP_PROBE=false
Environment=NEXT_TELEMETRY_DISABLED=1
ExecStartPre=/usr/bin/test -f /opt/auto-parse/.next/BUILD_ID
ExecStart=/usr/bin/npm run start -- -p 1007 -H 127.0.0.1
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillMode=control-group
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
